# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
loglevel=info

; ── FastAPI Backend ───────────────────────────────────────────────
[program:fastapi]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000 --log-level info
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/fastapi.log
stderr_logfile=/app/logs/fastapi.log
stdout_logfile_maxbytes=10MB
environment=PYTHONPATH="/app"

; ── Streamlit Frontend ────────────────────────────────────────────
[program:streamlit]
command=streamlit run frontend/app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true --server.fileWatcherType none
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/streamlit.log
stderr_logfile=/app/logs/streamlit.log
stdout_logfile_maxbytes=10MB
environment=PYTHONPATH="/app"